# emacs.d/server/emacs.service
[Unit]
Description=Emacs: the extensible, self-documenting text editor

[Service]
Type=forking
ExecStart=/usr/sbin/emacs --daemon
ExecStop=/usr/sbin/emacsclient --eval "(kill-emacs)"
Restart=always

[Install]
WantedBy=default.target
